# tps_so/tp6/ej4.cron

# min horas dias(mes) mes dias(semana)
  2 * * * * touch ~/cron_pruebas/archi_$(date +%H:%M).txt
